# nova-maintenance.timer
[Unit]
Description=NOVA Maintenance Timer (daily at 04:00)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
