# nova-maintenance.service
[Unit]
Description=NOVA Self-Maintenance
ConditionPathExists=/usr/local/lib/nova/scripts/maintenance.sh

[Service]
Type=oneshot
EnvironmentFile=-/etc/nova/nova.env
Environment=PATH=/home/ubuntu/.local/share/pnpm:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/lib/nova/scripts/maintenance.sh
StandardOutput=append:/var/log/nova/nova-maintenance.log
StandardError=append:/var/log/nova/nova-maintenance.log
