# nova-watchdog.service
[Unit]
Description=NOVA Emergency Watchdog
ConditionPathExists=/usr/local/lib/nova/scripts/watchdog.sh

[Service]
Type=oneshot
EnvironmentFile=-/etc/nova/nova.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/lib/nova/scripts/watchdog.sh
StandardOutput=append:/var/log/nova/nova-watchdog.log
StandardError=append:/var/log/nova/nova-watchdog.log
