# ops/nova-factory/systemd/nova-maintenance.service
[Unit]
Description=NOVA Self-Maintenance
ConditionPathExists=/usr/local/lib/nova/scripts/maintenance.sh

[Service]
Type=oneshot
EnvironmentFile=-/etc/nova/nova.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/lib/nova/scripts/maintenance.sh
NoNewPrivileges=true
UMask=0077
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/ubuntu/NOVA /home/ubuntu/backups /var/lib/nova /var/log/nova /tmp
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
StandardOutput=journal
StandardError=journal
