# ops/nova-factory/systemd/nova-dispatcher.service
[Unit]
Description=NOVA remote agent dispatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/nova/nova.env
ExecStart=/usr/local/lib/nova/nova-dispatcher.py dispatch-loop --limit 3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nova /var/log/nova /var/cache/nova
CapabilityBoundingSet=
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
