# ops/nova-factory/systemd/nova-orchestrator.timer
[Unit]
Description=Run NOVA orchestration coordination pass

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
RandomizedDelaySec=90s
Persistent=true
Unit=nova-orchestrator.service

[Install]
WantedBy=timers.target
